# eu.flawn.instantspaces.inject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>eu.flawn.instantspaces.inject</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>/Users/flawn/Documents/Projects/instantspaces/scripts/auto-inject.sh min0125</string>
    </array>

    <!-- Run when loaded (login) and keep alive (if it crashes) -->
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <!-- Ensure login session -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Environment PATH for lldb -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Optional: throttle restarts in case of rapid failures -->
    <key>ThrottleInterval</key>
    <integer>5</integer>
  </dict>
</plist>