# eu.flawn.instantspaces.inject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>eu.flawn.instantspaces.inject</string>

    <key>ProgramArguments</key>
    <array>
      <!-- TODO: update this to the absolute path of auto-inject.sh in your clone -->
      <string>/path/to/instantspaces/scripts/auto-inject.sh</string>
      <string>min0125</string>
    </array>

    <!-- Run when loaded (login) and only keep alive on failure (so success exits cleanly) -->
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <!-- Ensure login session -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Environment PATH for lldb -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Optional: throttle restarts in case of rapid failures -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Optional: log paths for debugging -->
    <key>StandardOutPath</key>
    <string>~/Library/Logs/instantspaces.auto-inject.out.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/instantspaces.auto-inject.err.log</string>
  </dict>
</plist>