# eu.flawn.instantspaces.inject.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>eu.flawn.instantspaces.inject</string>

    <!-- Installed automatically by: make install-agent -->
    <!-- WATCHER_PATH is replaced with /usr/local/bin/instantspaces-watcher by make install-agent -->
    <key>ProgramArguments</key>
    <array>
      <string>WATCHER_PATH</string>
    </array>

    <!-- Keep the watcher running always — it monitors Dock indefinitely -->
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>

    <!-- Ensure login session (required for Dock access) -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Environment PATH for lldb -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Log watcher output -->
    <key>StandardOutPath</key>
    <string>~/Library/Logs/instantspaces.watcher.out.log</string>
    <key>StandardErrorPath</key>
    <string>~/Library/Logs/instantspaces.watcher.err.log</string>
  </dict>
</plist>
